# data/watchdog.upstart.conf
description "Linux watchdog daemon"
author "Nathan Rennie-Waldock <nathan.renniewaldock@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

exec /usr/local/sbin/watchdog